# deploy/vito.service
[Unit]
Description=VITO Autonomous AI Agent
After=network.target postgresql.service
Wants=vito-watchdog.service

[Service]
Type=simple
User=vito
Group=vito
WorkingDirectory=/home/vito/vito-agent
EnvironmentFile=/home/vito/vito-agent/.env
ExecStart=/usr/bin/python3 -u main.py
Restart=always
RestartSec=10

# OOM Protection — systemd cgroup limits
# На сервере 8GB RAM. VITO получает макс 4GB.
MemoryMax=4G
MemoryHigh=3G

# CPU limit — не больше 2 ядер
CPUQuota=200%

# Лимит процессов — не плодить зомби
TasksMax=50

# Graceful shutdown via SIGTERM (main.py handles it)
KillSignal=SIGTERM
TimeoutStopSec=30

# Logging to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vito

[Install]
WantedBy=multi-user.target
